# system/postfix.service.d/restrictions.conf
[Service]
PrivateTmp=true
;DeviceAllow=closed

InaccessibleDirectories=/home /srv/http
ReadOnlyDirectories=/usr /etc
