# system/postfix.service.d/restrictions.conf
[Service]
PrivateTmp=true
PrivateDevices=true

ProtectHome=true
ProtectSystem=full
InaccessibleDirectories=/srv/http /var/lib/mysql
