# system/sshd.socket.d/network.conf
[Socket]
ListenStream=
ListenStream=144.76.107.12:512
ListenStream=[2a01:4f8:192:520b::2]:512
FreeBind=true
