# system/www-console@.service.d/restrictions.conf
[Service]
PrivateNetwork=false

InaccessibleDirectories=/srv/http/vhosts/admin.archlinux.de
InaccessibleDirectories=/srv/http/vhosts/forum.archlinux.de
InaccessibleDirectories=/srv/http/vhosts/pkgstats.archlinux.de
InaccessibleDirectories=/srv/http/vhosts/planet.archlinux.de
InaccessibleDirectories=/srv/http/vhosts/wiki.archlinux.de
