# system/openfire.service.d/restrictions.conf
[Service]
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true

ProtectHome=true
InaccessibleDirectories=/srv /var/lib/mysql
