# system/php-fpm@bbs.service.d/search.conf
[Service]
Environment="FPM_SOCKETS=/run/php-fpm/%i.socket=3,/run/php-fpm/%i-search.socket=4"
