# system/php-fpm7@bbs.service.d/search.conf
[Service]
Environment="FPM_SOCKETS=/run/php-fpm7/%i.socket=3,/run/php-fpm7/%i-search.socket=4"
