# system/dovecot.socket.d/network.conf
[Socket]
ListenStream=
ListenStream=0.0.0.0:143
ListenStream=[::]:143
ListenStream=0.0.0.0:993
ListenStream=[::]:993
