# system/php-fpm7@bbs.socket.d/search.conf
[Socket]
ListenStream=/run/php-fpm7/%i-search.socket
