# system/webhook.service.d/local.conf
[Service]
Type=simple
ExecStart=
ExecStart=/usr/bin/webhook -hooks /etc/webhook/hooks.json -verbose -ip 127.0.0.1
