# system/ntpd.service.d/restrictions.conf
[Service]
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true

ProtectHome=true
ProtectSystem=full
InaccessibleDirectories=/srv /var/lib/mysql
