# system/bbs-cleanup.timer
[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
