# system/forum-schedule-run.timer
[Timer]
OnUnitActiveSec=12min
OnBootSec=12min
RandomizedDelaySec=3min

[Install]
WantedBy=timers.target
