# system/nginx.socket
[Socket]
ListenStream=80
ListenStream=443
BindIPv6Only=ipv6-only
KeepAlive=true

[Install]
WantedBy=sockets.target
