# system/pierre-wp-cli-cron.timer
[Timer]
OnUnitActiveSec=27min
OnBootSec=27min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
