# system/planet-console@.service
[Service]
Type=oneshot

InaccessibleDirectories=/var/lib/mysql
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=full

User=php-planet
Group=php-planet

WorkingDirectory=/srv/http/vhosts/planet.archlinux.de/api
ExecStart=/srv/http/vhosts/planet.archlinux.de/api/bin/console %i

[Install]
WantedBy=multi-user.target
