# system/planet-console@app:update:feeds.timer
[Timer]
OnUnitActiveSec=30min
OnBootSec=30min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
