# system/spamassassin-update.timer
[Unit]
Description=daily spamassassin update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
