# system/spamassassin-update.service
[Unit]
Description=spamassassin update

[Service]
User=spamd
Group=spamd
Type=oneshot
ExecStart=-sa-update
ExecStart=-sa-compile
