# system/wiki-sitemaps.timer
[Timer]
OnUnitActiveSec=57min
OnBootSec=57min
RandomizedDelaySec=9min

[Install]
WantedBy=timers.target
