# system/wiki-run-jobs.timer
[Timer]
OnUnitActiveSec=42min
OnBootSec=42min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
