# system/www-console@.service
[Service]
Type=oneshot

User=php-www
Group=php-www

WorkingDirectory=/srv/http/vhosts/www.archlinux.de/api
ExecStart=/srv/http/vhosts/www.archlinux.de/api/bin/console %i

[Install]
WantedBy=multi-user.target
