# system/www-console@app:update:packages.timer
[Timer]
OnUnitActiveSec=2min
OnBootSec=2min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
