# system/www-console@app:update:mirror-popularities.timer
[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
