# system/www-console@app:update:releases.timer
[Timer]
OnUnitActiveSec=11min
OnBootSec=11min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
