# system/geoipupdate.service.d/path.conf
[Service]
ReadWritePaths=/usr/share/GeoIP
