# system/geoipupdate.service.d/restrictions.conf
[Service]
# Workround crash
SystemCallFilter=@resources
