# system/php-fpm@wiki.service.d/php-fpm-params.conf
[Service]
Environment="FPM_PARAMS=-d apc.shm_size=256M"
