# system/nginx.service.d/restrictions.conf
[Service]
PrivateTmp=true
;NoNewPrivileges=true
;DeviceAllow=closed

InaccessibleDirectories=/home /srv/vmail
ReadOnlyDirectories=/usr /etc
