# system/php-fpm@bbs.socket.d/search.conf
[Socket]
ListenStream=/run/php-fpm/%i-search.socket
