# system/spamassassin.service.d/update.conf
[Unit]
PartOf=spamassassin-update.service
